# scripts/pipeline/lofi-pipeline@.timer
[Unit]
Description=llooffiisounds — Daily Upload Timer (once per day at 06:00 UTC)

[Timer]
# Run once daily at 06:00 UTC
OnCalendar=*-*-* 06:00:00
# Catch up if the machine was off at the scheduled time
Persistent=true
# Add random delay up to 15 minutes to avoid exact-second hits
RandomizedDelaySec=900
Unit=lofi-pipeline@%i.service

[Install]
WantedBy=timers.target
